# etc/cron.hourly/oura_post
#!/bin/bash
/bin/python3 /etc/oura/oura_post_to_influxdb.py
