# crontab/oura_post
# Polls once per hour. Change as desired.
0 * * * * /usr/bin/python3 /etc/oura/oura_post_to_influxdb.py
